# container/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

trap -x
{
  term
  {
    foreground
    {
      import -u ! kill ${!}
    }
  }
}

exec nginx -g "daemon off;"
