# container/etc/services.d/delayed_job-log/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid danbooru tail -q -F -n0 /var/www/danbooru2/shared/log/production.log
